# common/supervisor/etc/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgresql]
user=postgres
command=/usr/lib/postgresql/9.5/bin/postgres -D "/opt/digitalbits/postgresql/data" -c config_file=/opt/digitalbits/postgresql/etc/postgresql.conf
stopsignal=INT
autostart=true
autorestart=true
priority=10

[program:digitalbits-core]
user=digitalbits
directory=/opt/digitalbits/core
command=/opt/digitalbits/core/bin/start
autostart=true
autorestart=true
priority=20

[program:frontier]
user=digitalbits
directory=/opt/digitalbits/frontier
command=/opt/digitalbits/frontier/bin/start
autostart=true
autorestart=true
priority=30
